# .config/systemd/user/cache.service
[Unit]
Description=Caching management

[Service]
Type=oneshot
ExecStart=/usr/bin/perl /home/enck/.local/bin/cache.pl hourly
