# .config/systemd/user/cache.timer
[Unit]
Description=Caching management

[Timer]
OnCalendar=hourly
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
