# .config/systemd/user/cache.service
[Unit]
Description=Caching management

[Service]
Type=oneshot
ExecStart=/usr/bin/drudge udb
ExecStart=/usr/bin/drudge udc
