# .config/systemd/user/cache.service
[Unit]
Description=Caching management

[Service]
Type=oneshot
ExecStart=/home/enck/.local/bin/localsvc.sh cache
