# .config/systemd/user/cache.service
[Unit]
Description=Caching management

[Service]
Type=oneshot
ExecStart=/usr/bin/drudge user.dir.backup
ExecStart=/usr/bin/drudge user.dir.cleanup
