# .config/systemd/user/html.service
[Unit]
Description=Open html pages

[Service]
Environment=DISPLAY=:0
KillMode=process
Type=oneshot
ExecStart=/home/enck/.bin/email webview
