# .config/systemd/user/html.path
[Unit]
Description=open firefox on html generation

[Path]
PathChanged=/home/enck/.tmp/html/

[Install]
WantedBy=default.target
