# .config/systemd/user/mail.service
[Unit]
Description=Local mail

[Service]
ExecStart=/usr/bin/perl /home/enck/.local/bin/mail.pl

[Install]
WantedBy=default.target
