# .config/systemd/user/mail.service
[Unit]
Description=Mail host service

[Service]
ExecStart=/home/enck/.local/bin/localsvc.sh mail_host

[Install]
WantedBy=default.target
