# .config/systemd/user/mailing.path
[Unit]
Description=trigger email on file request

[Path]
PathChanged=/home/enck/.cache/home/tmp/mail.trigger

[Install]
WantedBy=default.target
