# .config/systemd/user/maintain.timer
[Unit]
Description=Maintain user data (timer)

[Timer]
Unit=maintain.service
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
