# .config/systemd/user/maintain.service
[Unit]
Description=Maintaining user data

[Service]
Type=oneshot
ExecStart=/home/enck/.bin/maintain
