# .config/systemd/user/maintain.timer
[Unit]
Description=Maintain user data (timer)

[Timer]
Unit=maintain.service
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
