# .config/systemd/user/network.service
[Unit]
Description=network changing

[Service]
Type=oneshot
ExecStart=/home/enck/.bin/local-helpers network
