# .config/systemd/user/regenapps.service
[Unit]
Description=Regenerate app menu

[Service]
ExecStart=/usr/bin/perl /home/enck/.local/bin/cache.pl regen
