# .config/systemd/user/regenapps.service
[Unit]
Description=Regenerate app menu

[Service]
ExecStart=/home/enck/.local/bin/localsvc.sh regen
