# .config/systemd/user/regenapps.path
[Unit]
Description="Regenerate app menu"

[Path]
PathModified=/home/enck/.local/apps/

[Install]
WantedBy=default.target
