# .config/systemd/user/status.service
[Unit]
Description=Local status

[Service]
ExecStart=/home/enck/.local/bin/sys status

[Install]
WantedBy=default.target
