# .config/systemd/user/status.path
[Unit]
Description=update stats on request

[Path]
PathChanged=/home/enck/.tmp/.updatestats

[Install]
WantedBy=default.target
