# .config/systemd/user/sync.timer
[Unit]
Description=File syncing

[Timer]
Unit=sync.service
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
