# .config/systemd/user/sync.service
[Unit]
Description=File syncing

[Service]
Type=oneshot
ExecStart=/home/enck/.bin/syncing
