# .config/systemd/user/sync.service
[Unit]
Description=Syncing management

[Service]
ExecStart=/home/enck/.local/lib/localsvc.sh sync

[Install]
WantedBy=default.target
