# .config/systemd/user/sync.service
[Unit]
Description=Syncing management

[Service]
Type=oneshot
ExecStart=/home/enck/.local/lib/localsvc.sh sync
