# .config/systemd/user/sync.path
[Unit]
Description=change on wsw change

[Path]
PathChanged=/var/cache/wsw/

[Install]
WantedBy=default.target
