# .config/systemd/user/xsystem.service
[Unit]
Description=xsystem

[Service]
ExecStart=/home/enck/.local/bin/status

[Install]
WantedBy=multi-user.target
