# .config/systemd/user/xsystem.service
[Unit]
Description=xsystem

[Service]
Environment=DISPLAY=:0
ExecStart=/home/enck/.local/bin/xwindows system

[Install]
WantedBy=multi-user.target
