# aports/local/timeskew/timeskew.openrc
#!/sbin/openrc-run

command="/usr/bin/timeskew"
description="timeskew daemon"
pidfile="/run/timeskew.pid"
command_args=""

depend() {
	need ntp-server 
}
