# .config/go/wsw/wsw.service
[Unit]
Description=wpa_supplicant wrapper

[Service]
ExecStart=/usr/bin/wsw --mode daemon
Restart=always

[Install]
WantedBy=multi-user.target
