# .local/PKGBUILDs/localetc/etcgitstatus.path
[Unit]
Description=Check git status

[Path]
PathChanged=/var/cache/etcgit

[Install]
WantedBy=default.target
