# .local/PKGBUILDs/localetc/siteblock.service
[Unit]
Description=Domain blocking manager

[Service]
ExecStart=/home/enck/.local/bin/siteblock daemon

[Install]
WantedBy=multi-user.target
