# .local/PKGBUILDs/localetc/wsw.service
[Unit]
Description=wsw monitor

[Service]
ExecStart=/home/enck/.local/bin/wsw daemon

[Install]
WantedBy=multi-user.target
