# .config/systemd/etc/network-kill.service
[Unit]
Description=Perform network auth

[Service]
Type=oneshot
ExecStart=/home/enck/.config/systemd/etc/network exit
