# .config/systemd/etc/network-kill.path
[Unit]
Description=Setup network auth

[Path]
PathChanged=/home/enck/.cache/network

[Install]
WantedBy=default.target

