# .config/systemd/etc/network-auth.path
[Unit]
Description=Setup network auth

[Path]
PathChanged=/var/cache/networkill

[Install]
WantedBy=default.target

