# .config/systemd/etc/network-auth.path
[Unit]
Description=Setup network auth

[Path]
PathChanged=/home/enck/.cache/network.pid

[Install]
WantedBy=default.target

