# .config/systemd/etc/wsw-monitor.service
[Unit]
Description=wsw monitoring

[Service]
Type=oneshot
ExecStart=/home/enck/.bin/wsw disconnect
