# .config/systemd/etc/wsw-monitor.path
[Unit]
Description=wsw monitoring

[Path]
PathChanged=/var/cache/wsw/network

[Install]
WantedBy=default.target
