# .config/systemd/etc/wsw-monitor.service
[Unit]
Description=wsw monitoring

[Service]
Type=oneshot
ExecStart=/usr/bin/wsw disconnect
