# .config/systemd/etc/wsw.service
[Unit]
Description=wsw connections

[Service]
Type=oneshot
ExecStart=/usr/bin/wsw connect
