# .config/systemd/etc/wsw.path
[Unit]
Description=wsw connections

[Path]
PathChanged=/var/cache/wsw/mon.pid

[Install]
WantedBy=default.target
