# .bin/sys/rc.local
echo disable > /proc/acpi/ibm/bluetooth
exit 0
