# .bin/sys/sync.service
[Unit]
Description="File syncing"

[Service]
Type=oneshot
User=enck
ExecStart=/home/enck/.bin/sync
